# ops/systemd/dxm-digest.timer
[Unit]
Description=Run DXM Pulse weekly digest

[Timer]
OnCalendar=Mon 08:00
Persistent=true
Unit=dxm-digest.service

[Install]
WantedBy=timers.target
